# deploy/launchd/com.autotrader.papertrader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.autotrader.papertrader</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/tim/Development/auto-trader/.venv/bin/python</string>
        <string>-u</string>
        <string>scripts/paper_trader.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/tim/Development/auto-trader</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>60</integer>
    <key>StandardOutPath</key>
    <string>/Users/tim/Development/auto-trader/logs/paper_trader.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tim/Development/auto-trader/logs/paper_trader.log</string>
</dict>
</plist>
